# snowmelt-headless.service
[Unit]
Description=Snowmelt Control System (Headless)
Documentation=https://github.com/yourusername/snowmelt-control
After=network.target mosquitto.service
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/snowmelt_control
ExecStart=/home/pi/snowmelt_control/venv/bin/python /home/pi/snowmelt_control/main.py --no-gui
ExecStop=/bin/kill -SIGTERM $MAINPID
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=false
ProtectSystem=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target
