# snowmelt.service
[Unit]
Description=Snowmelt Control System
Documentation=https://github.com/JohnTheHobbit/Snowmelt-Control
After=network.target graphical.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/snowmelt_control
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
ExecStart=/home/pi/snowmelt_control/venv/bin/python /home/pi/snowmelt_control/main.py
ExecStop=/bin/kill -SIGTERM $MAINPID
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target